# docker/supervisord.conf
; ============================================================
;  Supervisor – mantém site + filas + scheduler + coleta vivos.
;  Workers reiniciam sozinhos se caírem (autorestart=true).
; ============================================================
[supervisord]
nodaemon=true
user=root
logfile=/var/www/html/storage/logs/supervisord.log
pidfile=/run/supervisord.pid

[program:php-fpm]
command=php-fpm -F
autostart=true
autorestart=true
priority=5
stdout_logfile=/var/www/html/storage/logs/php-fpm.log
stderr_logfile=/var/www/html/storage/logs/php-fpm.err.log

[program:nginx]
command=nginx -g 'daemon off;'
autostart=true
autorestart=true
priority=10
stdout_logfile=/var/www/html/storage/logs/nginx.log
stderr_logfile=/var/www/html/storage/logs/nginx.err.log

; ── Fila padrão: jobs gerais (push, sitemap, scores, links) ──
[program:queue-default]
command=php /var/www/html/artisan queue:work --queue=default --tries=3 --timeout=120 --sleep=3 --max-time=3600
autostart=true
autorestart=true
numprocs=1
user=www-data
stopwaitsecs=130
stdout_logfile=/var/www/html/storage/logs/services/queue-default.log
stderr_logfile=/var/www/html/storage/logs/services/queue-default.err.log

; ── Fila de coleta: scraping Python (jobs longos) ──
[program:queue-coleta]
command=php /var/www/html/artisan queue:work --queue=coleta --tries=2 --timeout=10800 --sleep=10 --max-time=14400
autostart=true
autorestart=true
numprocs=1
user=www-data
stopwaitsecs=20
stdout_logfile=/var/www/html/storage/logs/services/queue-coleta.log
stderr_logfile=/var/www/html/storage/logs/services/queue-coleta.err.log

; ── Scheduler: roda schedule:run a cada minuto em foreground ──
[program:scheduler]
command=php /var/www/html/artisan schedule:work
autostart=true
autorestart=true
user=www-data
stdout_logfile=/var/www/html/storage/logs/services/scheduler.log
stderr_logfile=/var/www/html/storage/logs/services/scheduler.err.log
